# com.snaptool.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.snaptool</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_PATH__</string>
        <string>__INSTALL_PATH__/snap_tool.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__INSTALL_PATH__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/snaptool.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/snaptool.err</string>
</dict>
</plist>
